# openlist/root/etc/s6-overlay/s6-rc.d/openlist/run
#!/usr/bin/with-contemv ash
# shellcheck shell=ash

exec /opt/openlist server --no-prefix
