# openlist/root/etc/s6-overlay/s6-rc.d/openlist/run
#!/usr/bin/with-contenv ash
# shellcheck shell=ash

exec s6-notifyoncheck -d -n 300 -w 1000 -c "nc -z localhost 5244" s6-setuidgid docker /opt/openlist/openlist server --no-prefix
