# openlist/root/etc/s6-overlay/s6-rc.d/aria2/run
#!/usr/bin/with-contemv ash
# shellcheck shell=ash

exec aria2c --conf-path /opt/aria2/config/aria2.conf
